# overlay/usr/lib/opkg/info/tor.list
/etc/init.d/tor
/etc/tor/torrc
/usr/sbin/tor
